# services/app-management/build/sysroot/usr/lib/systemd/system/6dainas-app-management.service
[Unit]
After=6dainas-message-bus.service
After=docker.service
Description=6DAiNAS-OS App Management Service

[Service]
ExecStartPre=/usr/bin/6dainas-app-management -v
ExecStart=/usr/bin/6dainas-app-management -c /etc/6dainas/app-management.conf
PIDFile=/var/run/6dainas/app-management.pid
Restart=always
Type=notify

[Install]
WantedBy=multi-user.target
